# curator4-1044393517.conf
[program:curator4-1044393517]
command=python3 /home/curator4/curator4-1044393517/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator4/curator4-1044393517/test.err.log
stdout_logfile=/home/curator4/curator4-1044393517/test.out.log